# sudo-apis/halogen-sudo-api.service
[Unit]
Description=Halogen Sudo API
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/home/msuser/halogen/sudo-apis
ExecStart=/bin/bash -c 'source /home/msuser/halogen/sudo-apis/venv/bin/activate && python3 /home/msuser/halogen/sudo-apis/main.py'
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=halogen-sudo-api

[Install]
WantedBy=multi-user.target
